# roles/fleet-server/files/fleet-yeet.timer
[Unit]
Description=Schedule for sending fleet data to loki

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target